# star-alliance-skills/skillsmith/assets/com.attax.skillsmith-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  skillsmith daily routine watchdog — macOS LaunchAgent.

  This LaunchAgent is the post-routine gate: it runs once a day at 08:00 local time,
  ~2 hours after the routine's expected ~06:00 finish window, and checks whether the
  skillsmith daily routine actually ran on schedule by reading
  star-alliance-skills/skillsmith/data/routine-heartbeat.json. If the heartbeat is
  stale or missing, it surfaces a clear, plain-English summary to the operator so a
  missed run cannot drift silently.

  PREFERRED mechanism is the native Claude Code routine (scheduled-tasks system, task
  `skillsmith-routine`, cron `0 6 * * *`) — manage it from the app's Scheduled panel.
  The watchdog above the routine is the optional app-closed fallback only; if you
  install this LaunchAgent, disable the native watchdog first so it doesn't
  double-run.

  PER-MACHINE: this plist is authored directly for this machine. All paths below are
  already absolute for /Users/atta — there is no REPLACE_ME templating; do not
  copy it to another user account without replacing the `/Users/atta` literal
  everywhere it appears. There is no `sed` install command because the file is
  meant to live on this box only.

  Install:    launchctl load ~/Library/LaunchAgents/com.attax.skillsmith-watchdog.plist
  Disable:    launchctl unload ~/Library/LaunchAgents/com.attax.skillsmith-watchdog.plist
  Retime:     edit StartCalendarInterval below, then unload + load again.
  Logs:       /Users/atta/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/watchdog-launchd.out
              /Users/atta/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/watchdog-launchd.err

  Installed by the skillsmith integration step (project-start / init). Do NOT run
  `launchctl load` by hand here — the integration step is responsible for both
  placing the file in ~/Library/LaunchAgents/ and loading it.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.attax.skillsmith-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/atta/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/scripts/routine_watchdog.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/atta/Documents/Claude/Projects/star-alliance</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/atta</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/atta/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/watchdog-launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/atta/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/watchdog-launchd.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
