# star-alliance-skills/skillsmith/assets/com.attax.skillsmith-routine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  skillsmith daily routine — macOS LaunchAgent.
  PREFERRED mechanism is the native Claude Code routine (scheduled-tasks system, task
  `skillsmith-routine`, cron `0 6 * * *`) — manage it from the app's Scheduled panel.
  This LaunchAgent is the optional app-closed fallback only; if you install it, disable the
  native routine first so it doesn't double-run.

  PER-MACHINE: the two absolute paths below are templated for the install user. Before
  copying, replace every `/Users/REPLACE_ME` with your own $HOME (run `echo $HOME`).
  Install:  sed "s#/Users/REPLACE_ME#$HOME#g" \
              "$HOME/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/assets/com.attax.skillsmith-routine.plist" \
              > ~/Library/LaunchAgents/com.attax.skillsmith-routine.plist
            launchctl load ~/Library/LaunchAgents/com.attax.skillsmith-routine.plist
  Disable:  launchctl unload ~/Library/LaunchAgents/com.attax.skillsmith-routine.plist
  Retime:   edit StartCalendarInterval below, then unload + load again.
  Logs:     $HOME/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/launchd.out
            (the run's own dated log is routine-logs/YYYY-MM-DD.log)
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.attax.skillsmith-routine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/REPLACE_ME/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/scripts/run_routine.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/Documents/Claude/Projects/star-alliance</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/REPLACE_ME</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/Documents/Claude/Projects/star-alliance/star-alliance-skills/skillsmith/routine-logs/launchd.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
